# profile.d/profile.conf.example
; Documentation: ../README.md (English) and ../README.pt-BR.md (Portuguese).
; PROFILE TEMPLATE: this file is not loaded because its extension is .example.
; Copy it to another *.conf file, set enabled = true, and replace the "dummy"
; values. Profile files are loaded in alphabetical order.
;
; [profile:name] — the name after "profile:" identifies the profile in logs
; and the HTTP response. It must be unique and descriptive.
[profile:example-dummy]
; enabled — false keeps this profile as a template; true enables processing.
enabled = false
; file — content destination. A relative path uses output_dir; use an absolute
; path, such as /absolute/path/dummy-output.log, outside the default directory.
file = dummy-output.log
; delivery — delivery method: file (default), fifo, or both. "both" writes the
; file and sends to the pipe without blocking the HTTP request.
delivery = both
; fifo_path — FIFO created automatically. A relative path uses output_dir.
fifo_path = pipes/dummy.fifo
; fifo_on_unavailable — for fifo/both: warn skips delivery if there is no
; reader; fail returns HTTP 503 so the sender can retry.
fifo_on_unavailable = warn
; rotate_max_bytes — rotates file/both delivery before the next message would
; exceed this size. The default is 10 MiB; set to 0 to explicitly disable it.
rotate_max_bytes = 10485760
; rotate_keep — number of rotated archive files to retain. The default is 10.
; Set to 0 to remove each archive immediately after creation.
rotate_keep = 10
; rotation_mode — the default is rename, which renames the active file and creates a
; new one. copytruncate copies to an archive then truncates the same inode;
; it helps legacy tail -f readers but can miss data they have not consumed.
rotation_mode = rename
; format — output format: raw, text, jsonl, or yaml. JSON is intentionally
; available only as JSON Lines, requiring a .jsonl or .ndjson output file.
format = text
; text_template — used only with format = text. It can use received JSON fields
; and {received_at}, {content_type}, {origin}, and {raw}.
text_template = [{received_at}] origin={origin} title={title} body={body}
; catch_all — true makes this a fallback profile, used only when no regular
; profile matches the message.
catch_all = false
; stop_after_match — true stops evaluating other profiles after this one is written.
stop_after_match = false
;
; origin — filters by the connection's exact source IP (optional).
origin = 192.0.2.10
; origin_cidr — filters by an IP range in CIDR notation (optional).
origin_cidr = 192.0.2.0/24
; origin_regex — filters the source text with a regular expression (optional).
origin_regex = ^192\\.0\\.2\\.[0-9]+$
;
; match.<field>.equals — requires an exact payload field match.
match.severity.equals = DUMMY_CRITICAL
; match.<field>.contains — requires the field to contain the supplied value.
match.type.contains = dummy-alarm
; match.<field>.regex — requires a regular expression match.
match.title.regex = ^Dummy.*$
; match.<field>.value — short form equivalent to exact equality.
match.source.value = dummy-source
; Nested fields are supported, including JSON serialized inside body.
match.body.alarmMetaData.severity.equals = DUMMY_WARNING
