# systemd/webhook-intake.service
[Unit]
Description=Webhook Intake receiver
Documentation=https://github.com/wuilber002/webhook-intake
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=whintake
Group=whintake
WorkingDirectory=@APP_DIR@
Environment=PYTHONUNBUFFERED=1
ExecStart=@APP_DIR@/.venv/bin/python @APP_DIR@/webhook.py --config @CONFIG_DIR@/config.ini
Restart=on-failure
RestartSec=5s
UMask=0077

# The service only needs to write received data under the state directory.
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=@STATE_DIR@
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
